# deploy/systemd/route53-update.timer
[Unit]
Description=Run Route53 update script every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
